# fundscreener.service
[Unit]
Description=Fund Screener Web Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/fundscreener
Environment="PATH=/var/www/fundscreener/venv/bin"
Environment="FLASK_CONFIG=production"
Environment="FLASK_PORT=8001"
ExecStart=/var/www/fundscreener/venv/bin/gunicorn --workers 4 --bind 127.0.0.1:8001 --timeout 120 wsgi:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
